# cron.weekly
#! /bin/bash

set -euo pipefail

srcdir="${PWD}"
workdir="$(mktemp -t -d "usage-bot-XXXXXX")"

./mkvenv "${workdir}/venv"

py="${workdir}/venv/bin/python3"

"${py}" usage-bot -wikitech -v -log
"${py}" usage-bot -osmwiki -v -log
